# lordfilm-bot.service
[Unit]
Description=LordFilm Telegram Bot
After=network.target

[Service]
Type=simple
User=kolla
Group=kolla
WorkingDirectory=/home/kolla/LordParser
ExecStart=/bin/bash -c 'cd /home/kolla/LordParser && /usr/bin/python3 -u bot.py'
Restart=on-failure
RestartSec=5s
Environment="PYTHONUNBUFFERED=1"
Environment="HOME=/home/kolla"
Environment="PYTHONPATH=/home/kolla/LordParser"
StandardOutput=file:/home/kolla/LordParser/bot.log
StandardError=file:/home/kolla/LordParser/bot.log

[Install]
WantedBy=multi-user.target
